# factorio-idle.timer
[Unit]
Description=Factorio zero player server shutdown timer

[Timer]
OnActiveSec=10s
OnUnitActiveSec=5m

[Install]
WantedBy=timers.target
