# factorio.service
[Unit]
Description=Factorio headless server

[Service]
User=factorio
ExecStart=/factorio/run.sh

[Install]
WantedBy=multi-user.target
